# deploy/dnsd.service
[Unit]
Description=dnsd — host DNS policy resolver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/dnsd/dnsd.env
ExecStart=/usr/local/bin/dnsd \
  --listen ${DNSD_LISTEN:-127.0.0.1:51920} \
  --token ${DNSD_TOKEN:-dev-token} \
  --dns-listen ${DNSD_DNS_LISTEN:-127.0.0.1:5353} \
  --upstream ${DNSD_UPSTREAM:-1.1.1.1:53,8.8.8.8:53}
Restart=on-failure
RestartSec=2
# AmbientCapabilities=CAP_NET_BIND_SERVICE
# User=dnsd

[Install]
WantedBy=multi-user.target
